# util/docker/web/runit/sftpgo/run
#!/bin/bash

echo 'Spinning up SFTP process...'

cd /var/playcast/sftpgo

sudo -E -u playcast sftpgo --config-dir=/var/playcast/sftpgo serve -l "" > /proc/1/fd/1 2> /proc/1/fd/2
